# debian/almanach-api.prerm
#!/bin/sh

set -e

invoke-rc.d almanach-api stop

exit 0
