# debian/almanach-collector.prerm
#!/bin/sh

set -e

invoke-rc.d almanach-collector stop

exit 0
